# packaging/deb/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade)
        if command -v systemctl >/dev/null 2>&1; then
            systemctl stop gatel.service 2>/dev/null || true
            systemctl disable gatel.service 2>/dev/null || true
        fi
        ;;
esac

exit 0
